# web/deploy/rt-dashboard.service
[Unit]
Description=Rotten Tomatoes Analytics Dashboard
After=network.target

[Service]
Type=exec
User=jakelehner
Group=jakelehner
WorkingDirectory=/home/jakelehner/rotten-tomatoes-analysis/web
ExecStart=/home/jakelehner/.local/bin/uv run uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 1
StandardOutput=append:/home/jakelehner/rotten-tomatoes-analysis/dashboard.log
StandardError=append:/home/jakelehner/rotten-tomatoes-analysis/dashboard.log

# Environment variables for the dashboard config
Environment=RT_DB_PATH=/home/jakelehner/rotten-tomatoes-analysis/reviews.db
Environment=RT_MOVIES_JSON=/home/jakelehner/rotten-tomatoes-analysis/movies.json
Environment=RT_HOST=127.0.0.1
Environment=RT_PORT=8000

# Memory limits to protect the scraper's headroom
MemoryHigh=150M
MemoryMax=200M

# Restart on crash, but not on clean stop
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
